# dropbox/.config/systemd/user/dropbox.service
[Unit]
Description=Dropbox

[Service]
ConditionPathExists=${HOME}/Dropbox
ExecStart=/usr/bin/env ${HOME}/.dropbox-dist/dropboxd
StandardOutput=null
StandardError=null
SyslogIdentifier=dropbox
Restart=always

[Install]
WantedBy=default.target
